# roles/php_app/templates/rc.local.j2
{% for item in extra_directories %}
mkdir {{ item }}
sudo chmod 777 {{ item }}
{% endfor %}

exit 0